# SysAdmin/crontab/ADD_crontab.sh
#!/bin/bash
crontab -e


# IN PROGRESS...
#USER="-e"
#if (( $1 ! 0 )) ; then
#	crontab -u $1
#elif
#	crontab $USER
