# SysAdmin/crontab/STOP_crontab.sh
#!/bin/bash
service cron status && \
service cron stop && \
service cron status
